# systemd/sway.service
[Unit]
Description=Sway Wayland Compositor

[Service]
ExecStart=/home/zinnia/sway

[Install]
WantedBy=default.target
